# etc/rc.openbsd
#!/bin/ksh
#
# $OpenBSD: flamingzombies, ${VERSION} ${TIMESTAMP} rpe Exp $

daemon="/opt/flamingzombies-${VERSION}/bin/fz"

. /etc/rc.d/rc.subr

rc_start() {
        ${daemon} ${daemon_flags} &
}

rc_configtest() {
        ${daemon} ${daemon_flags} -n || return 1
}

rc_cmd $1
